# state/base/files/rc.local
# Customized rc.local
# Set deadline scheduler by default
echo deadline > /sys/block/sda/queue/scheduler
